# dot_config/systemd/user/cronpkg.service
[Unit]
Description=CronPkg - Sync package lists to GitHub

[Service]
Type=oneshot
ExecStart=/home/joerg/Documents/pkglist/cronpkg --run
