# dot_config/systemd/user/pkgsync.service
[Unit]
Description=PkgSync - Sync package lists to GitHub

[Service]
Type=oneshot
ExecStart=/home/joerg/.local/bin/pkgsync --run
